# docs/node-agent-install/assets/systemd/nyxid-node.service
[Unit]
Description=NyxID Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{RUN_AS_USER}}
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart={{BINARY_PATH}} --log-level {{LOG_LEVEL}} start --config {{CONFIG_DIR}}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
